# botrunner.service
[Unit]
Description=CyberBot
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /home/age_2023/Saidazimxon/cyber-bot/main.py
WorkingDirectory=/home/age_2023/Saidazimxon/cyber-bot
StandardOutput=journal
StandardError=journal
Restart=on-failure

[Install]
WantedBy=multi-user.target
